# deploy/com.mv4.preopen.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd PRE-OPEN agent. Install:
       cp deploy/com.mv4.preopen.plist ~/Library/LaunchAgents/
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.mv4.preopen.plist
     Verify: launchctl list | grep mv4   ·   logs: /tmp/mv4_preopen.log
     Fires at 07:30 / 08:00 / 08:30 LOCAL, Mon-Fri (brackets the ET pre-open window whether the host is on
     Eastern, Central, or Mountain). The script gates the actual trade on the BROKER CLOCK (market closed,
     next open today, <=100 min away) and a once-per-ET-day lock, so exactly one in-window run trades and
     the rest no-op — robust to a host-timezone change. Reruns on wake if asleep. Rebalances the CALENDAR
     books (weekly/biweekly/monthly) AT THE OPEN via limit-on-open orders; each book self-gates to its own
     first-of-period day. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.mv4.preopen</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/david/PyCharmMiscProject/ai-paper-trader/deploy/preopen_rebalance.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/tmp/mv4_preopen.log</string>
  <key>StandardErrorPath</key><string>/tmp/mv4_preopen.log</string>
  <key>SoftResourceLimits</key>
  <dict><key>NumberOfFiles</key><integer>10240</integer></dict>
  <key>HardResourceLimits</key>
  <dict><key>NumberOfFiles</key><integer>10240</integer></dict>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
