# deploy/com.mv4.preclose.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd PRE-CLOSE trade agent. Install:
       cp deploy/com.mv4.preclose.plist ~/Library/LaunchAgents/
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.mv4.preclose.plist
     Verify: launchctl list | grep mv4   ·   logs: /tmp/mv4_preclose.log
     Runs 14:40 local Central (= 15:40 ET, 20 min BEFORE the close), Mon-Fri; reruns on wake if asleep. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.mv4.preclose</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/david/PyCharmMiscProject/ai-paper-trader/deploy/preclose_trade.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>40</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/tmp/mv4_preclose.log</string>
  <key>StandardErrorPath</key><string>/tmp/mv4_preclose.log</string>
  <key>SoftResourceLimits</key>
  <dict><key>NumberOfFiles</key><integer>10240</integer></dict>
  <key>HardResourceLimits</key>
  <dict><key>NumberOfFiles</key><integer>10240</integer></dict>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
