# deploy/com.mv4.midday.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd MIDDAY agent (informational, no trading). Install:
       cp deploy/com.mv4.midday.plist ~/Library/LaunchAgents/
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.mv4.midday.plist
     Verify: launchctl list | grep mv4   ·   logs: /tmp/mv4_midday.log
     Runs 11:30 local Central (= 12:30 ET), Mon-Fri; reruns on wake if the machine was asleep. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.mv4.midday</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/david/PyCharmMiscProject/ai-paper-trader/deploy/midday_summary.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/tmp/mv4_midday.log</string>
  <key>StandardErrorPath</key><string>/tmp/mv4_midday.log</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
